# terraform/k3s/ip.service
[Unit]
Description=Run my task every 5 seconds

[Service]
Type=simple
ExecStart=/tmp/update_ip.sh
