# terraform/k3s/ip.timer
[Unit]
Description=Runs ip.service every 5 seconds

[Timer]
OnBootSec=5sec
OnUnitActiveSec=5sec
Unit=ip.service

[Install]
WantedBy=timers.target